# deploy/supervisord_mynas.ini
[program:gunicorn]
command=/data1/data/myNAS/venv/bin/gunicorn -w 8 -b 0.0.0.0:4999 mynasflask:app
directory=/data1/data/myNAS
user=app
autostart=true
autorestart=true
redirect_stderr=true
startretries=3
stderr_logfile=/data1/data/myNAS/err.log
stdout_logfile=/data1/data/myNAS/out.log